# crontab
*/5 * * * * cd /var/www/html && su -s /bin/sh www-data -c 'php xcat Cron'
